# where_ya_at_bulb/rootfs/etc/services.d/where_ya_at_bulb/run
#!/usr/bin/with-contenv bashio

cd /app
exec /usr/bin/python3 -m server
